# nibble2-win32-x64/resources/app/nibble-server.service
[Unit]
Description=Backend for Nibble

[Service]
RestartSec=5
Restart=always
ExecStart=/usr/bin/npm run --prefix $NIBBLE_DIR dev

[Install]
WantedBy=multi-user.target

